# puppet/modules/common/files/etc/cron.d/nssdb
@reboot root /var/db/nssdb_update.sh >/dev/null 2>/dev/null
45 * * * * root /var/db/nssdb_update.sh >/dev/null 2>/dev/null
