# puppet/modules/common/files/etc/cron.d/puppet
@reboot root /usr/site/rcac/sbin/run_puppet 2>/dev/null
30 * * * * root /usr/site/rcac/sbin/run_puppet 2>/dev/null
