# puppet/modules/common/files/etc/cron.d/syncfiles
# This cron job tries to emulate xCAT's syncfiles functionality for the ADM
@reboot root /usr/site/rcac/sbin/syncfiles 2>/dev/null
30 * * * * root /usr/site/rcac/sbin/syncfiles 2>/dev/null
