# puppet/modules/common/files/etc/profile.d/000_rcac_prompt.csh
# Set (t)csh Prompt

# For interactive shells only.
if ( $?prompt ) then
	set promptchars = '$#'
        set prompt = '%n@%m:%B%~%b %# '
endif
