# puppet/modules/common/files/etc/profile.d/hpss.csh
# HPSS initialization script (csh, tcsh)

set path = ( $path /opt/hpss/bin /opt/hsi/bin )

